# synology-drive-client-bin/synology-drive.service
[Unit]
Description=Synology Drive Client Daemon
After=network.target

[Service]
Type=simple
ExecStart=/opt/Synology/SynologyDrive/package/cloudstation/bin/cloud-drive-daemon
WorkingDirectory=/opt/Synology/SynologyDrive/package/cloudstation/bin
Restart=on-failure

[Install]
WantedBy=default.target
