# synology-drive-client-bin/synology-drive.service
[Unit]
Description=Synology Drive Client
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/synology-drive start
Environment="QT_QPA_PLATFORM=xcb"
Restart=on-failure

[Install]
WantedBy=default.target
